# modules/sandpit/files/sandpit_app.conf
start on startup
console log

description "Sandpit app start script"
author "Tim Moran <tim@t-moran.com>"

script
  exec /vagrant/run.sh
  respawn
end script