# puppet/modules/sandpit/files/sandpit_celery.conf
start on startup
console log

description "Sandpit celery start script"
author "Tim Moran <tim@t-moran.com>"

script
  exec /vagrant/run_celery.sh
  respawn
end script